# private_dot_config/systemd/user/ssh-agent.service.d/override.conf
[Service]
Environment=IDENTITY_LIFETIME=8h
RuntimeMaxSec=infinity
